# Primary Cubesat/PI0/rc.local
#!/bin/sh -e

python3 /home/pi/masterCameraCode.py & python3 /home/pi/music.py &

exit 0
