# Primary Cubesat/rc.local
#!/bin/sh -e

sudo python3 /home/pi/main.py 2> error.txt &

exit 0
